# docker/config/supervisord-apache2.conf
[program:apache2]
command=apachectl -D "FOREGROUND" -k start
numprocs=1
autostart=true
autorestart=true
